# mongodb.cron
* 3 * * * root cd /var/www && sh mongodb.dump.sh
